# systemd/recover-openclaw-agents.service
[Unit]
Description=Recover OpenClaw + Claude Code agents after reboot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=forge
WorkingDirectory=/home/forge/.openclaw/workspace
ExecStart=/home/forge/.openclaw/workspace/skills/gsd-code-skill/scripts/recover-openclaw-agents.sh --registry /home/forge/.openclaw/workspace/skills/gsd-code-skill/config/recovery-registry.json

[Install]
WantedBy=multi-user.target
